# board/raspberrypi/fs_commnode/etc/systemd/system/macsec.service
[Unit]
Description=macsec networking
After=multi-user.target

[Service]
Type=oneshot
WorkingDirectory=/opt/nk-macsec
ExecStart=/opt/nk-macsec/macsec.sh
RemainAfterExit=true
